# Restaurantes/despliegue/supervisor.conf
[program:gunicorn]
command=/usr/local/bin/gunicorn Restaurantes.wsgi  --bind 0.0.0.0:8000
directory=/var/Volumes/manage.py
user=ubuntu
autostart=true
autorestart=true
redirect_stderr=true
